# bin/99-power-events.rules
SUBSYSTEMS=="power_supply", ACTION=="change", RUN+="/usr/local/bin/power_events"
